# nbnc.service
[Unit]
Description=null (transparent) bnc

[Service]
ExecStart=/usr/bin/nbnc /etc/nbnc.conf -
Restart=always

[Install]
WantedBy=multi-user.target
